# systemd/timers/lightdmbg.timer
[Unit]
Description=Fetch a random 4k UHD wallpaper from unsplash every 5 minutes

[Timer]
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
