# systemd/timers/lightdmbg.service
[Unit]
Description=background from unsplash

[Service]
Type=oneshot
ExecStart=cp /tmp/wallpaper.jpg /usr/share/backgrounds/wallpaper.jpg
